# scripts/podman/idlehands.container.in
# IdleHands gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into idlehands's ~/.config/containers/systemd/
# {{IDLEHANDS_HOME}} is replaced at install time.

[Unit]
Description=IdleHands gateway (rootless Podman)

[Container]
Image=idlehands:local
ContainerName=idlehands
UserNS=keep-id
Volume={{IDLEHANDS_HOME}}/.idlehands:/home/node/.idlehands
EnvironmentFile={{IDLEHANDS_HOME}}/.idlehands/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
