# contrib/idlehands-bot.service
# Idle Hands Telegram Bot — systemd user service
# Install: cp contrib/idlehands-bot.service ~/.config/systemd/user/
# Usage:
#   systemctl --user daemon-reload
#   systemctl --user enable --now idlehands-bot
#   journalctl --user -u idlehands-bot -f

[Unit]
Description=Idle Hands Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.npm-global/bin/idlehands bot telegram
Restart=on-failure
RestartSec=5

# Required: set your bot token here or in config.json
Environment=IDLEHANDS_TG_TOKEN=your-token-here

# Optional: override endpoint (default from config.json)
# Environment=IDLEHANDS_ENDPOINT=http://192.168.68.118:8080/v1

# Optional: override allowed users (comma-separated IDs)
# Environment=IDLEHANDS_TG_ALLOWED_USERS=123456789

# Optional: override working directory
# Environment=IDLEHANDS_TG_DIR=~/projects

[Install]
WantedBy=default.target
